# ansible/roles/app/files/fastapi-app.socket
[Unit]
Description=FastAPI App Socket

[Socket]
ListenStream=/run/fastapi-app.sock

[Install]
WantedBy=sockets.target
